# debian/mx3.postrm
#!/bin/sh
set -e

case "$1" in
    purge|remove)
        systemctl disable --now mx3 2>/dev/null || true
    ;;
esac

#DEBHELPER#
